# scripts/install/kairix-fetch-secrets.service.example
# systemd oneshot that hydrates /run/secrets/kairix.env from Azure Key Vault.
#
# Copy to /etc/systemd/system/kairix-fetch-secrets.service, edit
# placeholders if needed, then:
#
#   sudo systemctl daemon-reload
#   sudo systemctl enable kairix-fetch-secrets.service
#
# This unit is `Required=` by kairix.service (see kairix.service.example),
# so kairix never starts before the tmpfs secrets file has been written.
#
# Why this file ships with kairix:
# /run is tmpfs — secrets vanish on every reboot. Without this oneshot
# pinned to boot, kairix.service comes up against a missing
# /run/secrets/kairix.env and silently degrades to BM25-only search
# (#167). Enabling this oneshot is the durable fix.

[Unit]
Description=Kairix — fetch Azure Key Vault secrets to /run/secrets/kairix.env
After=network-online.target
Wants=network-online.target
Before=kairix.service

# Don't run on hosts with no Key Vault config — kairix can be deployed
# with secrets baked into /opt/kairix/secrets.env instead.
ConditionPathExists=/opt/kairix/bin/fetch-secrets.sh

[Service]
Type=oneshot
RemainAfterExit=yes

# The fetch script is shipped with kairix (see scripts/install/fetch-secrets.sh).
# It uses the VM's managed identity OR an explicit service principal to pull:
#   kairix-llm-api-key, kairix-llm-endpoint, kairix-embed-api-key,
#   kairix-embed-endpoint, kairix-neo4j-password
# from the configured Key Vault and writes them as KEY=VALUE lines to
# /run/secrets/kairix.env (mode 0640, owner root:kairix).
ExecStart=/opt/kairix/bin/fetch-secrets.sh

# If the fetch fails (network, KV permissions), surface the failure
# instead of leaving an empty secrets file in place.
ExecStartPost=/usr/bin/test -s /run/secrets/kairix.env

# Run as root because /run/secrets needs root mode bits and tmpfs
# permissions; the script chowns the result to kairix:kairix.
User=root
Group=root

[Install]
WantedBy=multi-user.target
